# files/etc/systemd/system/scanner-temp-watch.timer
[Unit]
Description=Run the scanner SoC temperature watcher every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
# Catch up if the Pi was asleep/off when a run was due.
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
